# rootfs/etc/s6-overlay/s6-rc.d/razorfy/run
#!/command/execlineb -P
exec s6-setuidgid abc /app/razorfy.pl